# systemd/wireguard-route.service
# /etc/systemd/system/wireguard-route.service
[Unit]
Description=WireGuard route for VPN clients
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wireguard-route.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
